# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=Set desktop wallpaper

[Service]
Type=simple
ExecStart=/home/jacob/bin/rotate-wallpaper

[Install]
WantedBy=default.target
