# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change desktop wallpaper every 10 minutes

[Timer]
OnActiveSec=5
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
