# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change desktop wallpaper every 10 minutes

[Timer]
OnActiveSec=20
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
